# files/etc/profile.d/calligra.sh
#
# calligra.sh
#
# Environment modifications for Calligra
#
pathappend  /opt/calligra/bin   PATH
pathappend  /opt/calligra/share XDG_DATA_DIRS
pathappend  /opt/calligra/share/calligraactive/qml  QML_IMPORT_PATH
pathappend  /opt/calligra/share/calligraactive/qml  QML2_IMPORT_PATH
#
# I THINK that should do it!
# End calligra.sh
