# files/etc/profile.d/enlightenment.sh
#!/bin/bash
#
# Begin enlightenment.sh
#
pathappend /opt/enlightenment/bin               PATH
pathappend /opt/enlightenment/lib/pkgconfig     PKG_CONFIG_PATH  
pathappend /opt/enlightenment/share/man         MANPATH  
pathappend /opt/enlightenment/share             XDG_DATA_DIRS
pathappend /opt/etc/xdg/enlightenment           XDG_CONFIG_DIRS

# End enlightenment.sh

