# files/etc/profile.d/pyqt5.sh
#
# Begin /etc/profile.d/pyqt5.sh

PYQT5_SIP_DIR=/usr/share/sip/PyQt5

pathappend /opt/pyqt5/bin 	PATH
pathappend /opt/pyqt5 		PYTHONPATH

export PYQT5_SIP_DIR

# End /etc/profile.d/pyqt5.sh


